# roles/nodejs/templates/systemd.service
[Unit]
Description = {{COMPONENT}} Service
After=network.target

[Service]
User={{USER}}
ExecStart=/bin/npm start
SyslogIdentifier={{COMPONENT}}
WorkingDirectory=/home/{{USER}}/{{COMPONENT}}/
# StandardOutput=syslog
# StandardError=syslog
# Environment=MONGO_URL="mongodb://mongodb-{{ENV}}.roboshop.internal:27017/catalogue"

[Install]
WantedBy=multi-user.target